# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:acrel]
command=python /app/acrel.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/acrel.log
stderr_logfile=/var/log/supervisor/acrel.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:pressure]
command=python /app/pressure.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/pressure.log
stderr_logfile=/var/log/supervisor/pressure.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:mqtt_to_mysql]
command=python /app/mqtt_to_mysql.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mqtt_to_mysql.log
stderr_logfile=/var/log/supervisor/mqtt_to_mysql.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

[program:dooble]
command=python /app/dooble.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dooble.log
stderr_logfile=/var/log/supervisor/dooble.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
